# supervisord.conf
# =============================================================================
# supervisord.conf — CiteRag Process Manager Configuration
# =============================================================================
# supervisord runs both FastAPI and Streamlit inside the same Docker container.
# It automatically restarts either service if it crashes.
#
# Logs go to stdout/stderr so Docker captures them with `docker logs`.
# =============================================================================

[supervisord]
nodaemon=true                    ; run in foreground (required for Docker)
logfile=/dev/null                ; no supervisord log file — Docker captures stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

# ── FastAPI Backend ───────────────────────────────────────────────────────────
[program:fastapi]
command=uvicorn main:app --host 0.0.0.0 --port 8000 --workers 1
directory=/app                   ; working dir must match where main.py lives
autostart=true                   ; start on supervisord boot
autorestart=true                 ; restart if it crashes
startretries=3                   ; retry up to 3 times before giving up
stdout_logfile=/dev/stdout       ; pipe to Docker logs
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── Streamlit Frontend ────────────────────────────────────────────────────────
[program:streamlit]
command=streamlit run frontend.py
    --server.port=7860
    --server.address=0.0.0.0
    --server.headless=true
    --server.enableCORS=false
    --server.enableXsrfProtection=false
directory=/app                   ; working dir must match where frontend.py lives
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
